# deploy/systemd/ytd_web.service
[Unit]
Description=Web YTD FastAPI service
After=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__VENV_DIR__/bin/python __APP_DIR__/web_ytd.py
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
